# fabric/lib/templates/mount.automount
[Unit]
Description=Automount {{ source }}
[Automount]
Where={{ dest }}
[Install]
WantedBy=multi-user.target
WantedBy=docker.service
